# webhook-listener.service
[Unit]
Description=YachtThot GitHub Webhook Listener
After=network.target

[Service]
Type=simple
User=thefleshgordon
WorkingDirectory=/home/thefleshgordon/YachtThot-v2
ExecStart=/home/thefleshgordon/YachtThot-v2/venv/bin/gunicorn -w 1 -b 0.0.0.0:5000 webhook_listener:app
Restart=always
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target